# chatgpt-guard.socket
[Socket]
ListenStream=80

[Install]
WantedBy=graphical.target
