# chatgpt-guard-ssl.service
[Service]
ExecStart=python %Y/redirect.py --ssl
User=student
WorkingDirectory=%Y

[Install]
WantedBy=graphical.target
