# chatgpt-guard-ssl.socket
[Socket]
ListenStream=443

[Install]
Also=chatgpt-guard-ssl.socket
