# templates/ai-guard-ssl.socket.j2
[Socket]
ListenStream={{ ai_redirect_localhost_ip }}:443

[Install]
Also=ai-guard-ssl.socket
