# templates/ai-guard.socket.j2
[Socket]
ListenStream={{ ai_redirect_localhost_ip }}:80

[Install]
WantedBy=graphical.target
