# templates/ai-guard.service.j2
[Service]
ExecStart=python redirect.py
WorkingDirectory={{ ai_redirect_script_folder }}

[Install]
Also=ai-guard.socket
